# yt6801-reinstall.service
[Unit]
Description=Install or reinstall YT6801 driver if not loaded
After=network-pre.target

[Service]
Type=oneshot
WorkingDirectory=/opt/yt6801-auto-installer
ExecStart=/opt/yt6801-auto-installer/install_yt6801_if_needed.sh
ExecStartPost=/opt/yt6801-auto-installer/check_yt6801_and_reboot.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
